# ansible/templates/ring-server.conf
exec {{project_dir}}/scripts/ring-server.sh
start on startup
chdir {{project_dir}}
