# ansible/templates/ring-server.conf
exec {{script_dir}}/ring-server.sh
start on startup
chdir {{project_dir}}
